# net-wireless/dump1090/files/dump1090.initd-r1
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

name="dump1090 Flight Aware Daemon"
command="/usr/bin/dump1090"
command_args="--quiet --net --enable-agc --write-json /run/dump1090-fa/"
command_user="dump1090:dump1090"
pidfile="/run/dump1090-fa.pid"
command_background=true

depend() {
	need lighttpd
}
